# config/worker.ini
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:product-synchronization]
command=php /app/bin/console messenger:consume product_synchronization -vv
numprocs=2
startsecs=0
autostart=true
autorestart=true
startretries=10
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true